# deploy/systemd/partijgedrag-serve.service
[Unit]
Description=Partijgedrag web and API server
Documentation=https://codeberg.org/partijgedrag/partijgedrag
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=simple
User=partijgedrag
Group=partijgedrag
WorkingDirectory=/opt/partijgedrag
EnvironmentFile=/etc/partijgedrag/partijgedrag.env
Environment=PATH=/opt/partijgedrag/bin:/usr/local/bin:/usr/bin
# The sync timer unit handles syncing in this setup; disable the built-in scheduler.
Environment=SYNC_INTERVAL=0
ExecStart=/usr/bin/env partijgedrag serve
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/partijgedrag /var/log/partijgedrag

[Install]
WantedBy=multi-user.target
