# etl/crontab
0 * * * * /app/etl -after yesterday >> /var/log/etl.log 2>&1